# core/roles/02_vmhosts/templates/etc/systemd/system/ipxehttpd.service
[Service]
ExecStart=busybox httpd -p {{ net_v4 }}:8000 -u nobody:nogroup -h /srv/ipxe -f 
Restart=always
RestartSec=30